# grubdev/cde/optional/open-vm-tools.tcz.list
/usr/local/share/open-vm-tools/messages/fr/toolboxcmd.vmsg
/usr/local/share/open-vm-tools/messages/fr/vmtoolsd.vmsg
/usr/local/share/open-vm-tools/messages/zh_TW/toolboxcmd.vmsg
/usr/local/share/open-vm-tools/messages/zh_TW/vmtoolsd.vmsg
/usr/local/share/open-vm-tools/messages/es/toolboxcmd.vmsg
/usr/local/share/open-vm-tools/messages/es/vmtoolsd.vmsg
/usr/local/share/open-vm-tools/messages/en/toolboxcmd.vmsg
/usr/local/share/open-vm-tools/messages/en/vmtoolsd.vmsg
/usr/local/share/open-vm-tools/messages/ja/toolboxcmd.vmsg
/usr/local/share/open-vm-tools/messages/ja/vmtoolsd.vmsg
/usr/local/share/open-vm-tools/messages/ko/toolboxcmd.vmsg
/usr/local/share/open-vm-tools/messages/ko/vmtoolsd.vmsg
/usr/local/share/open-vm-tools/messages/de/toolboxcmd.vmsg
/usr/local/share/open-vm-tools/messages/de/vmtoolsd.vmsg
/usr/local/share/open-vm-tools/messages/zh_CN/toolboxcmd.vmsg
/usr/local/share/open-vm-tools/messages/zh_CN/vmtoolsd.vmsg
/usr/local/share/open-vm-tools/messages/it/toolboxcmd.vmsg
/usr/local/share/open-vm-tools/messages/it/vmtoolsd.vmsg
/usr/local/etc/init.d/open-vm-tools
/usr/local/etc/udev/rules.d/99-vmware-scsi-udev.rules
/usr/local/tce.installed/open-vm-tools
/usr/local/sbin/mount.vmhgfs
/usr/local/lib/libvmtools.so.0.0.0
/usr/local/lib/libguestlib.so
/usr/local/lib/libguestlib.so.0
/usr/local/lib/libguestlib.so.0.0.0
/usr/local/lib/libappmonitor.so.0
/usr/local/lib/open-vm-tools/plugins/common/libvix.so
/usr/local/lib/open-vm-tools/plugins/common/libhgfsServer.so
/usr/local/lib/open-vm-tools/plugins/vmsvc/libserviceDiscovery.so
/usr/local/lib/open-vm-tools/plugins/vmsvc/libtimeSync.so
/usr/local/lib/open-vm-tools/plugins/vmsvc/libpowerOps.so
/usr/local/lib/open-vm-tools/plugins/vmsvc/libvmbackup.so
/usr/local/lib/open-vm-tools/plugins/vmsvc/libguestInfo.so
/usr/local/lib/open-vm-tools/plugins/vmsvc/libappInfo.so
/usr/local/lib/open-vm-tools/serviceDiscovery/scripts/get-listening-process-perf-metrics.sh
/usr/local/lib/open-vm-tools/serviceDiscovery/scripts/get-versions.sh
/usr/local/lib/open-vm-tools/serviceDiscovery/scripts/get-listening-process-info.sh
/usr/local/lib/open-vm-tools/serviceDiscovery/scripts/get-connection-info.sh
/usr/local/lib/libhgfs.so.0.0.0
/usr/local/lib/libvmtools.so.0
/usr/local/lib/libhgfs.so.0
/usr/local/lib/libhgfs.so
/usr/local/lib/libappmonitor.so
/usr/local/lib/libvmtools.so
/usr/local/lib/libappmonitor.so.0.0.0
/usr/local/bin/vmware-vmblock-fuse
/usr/local/bin/vmware-checkvm
/usr/local/bin/vmhgfs-fuse
/usr/local/bin/vmware-xferlogs
/usr/local/bin/vmware-rpctool
/usr/local/bin/vmware-namespace-cmd
/usr/local/bin/vmware-hgfsclient
/usr/local/bin/vmware-toolbox-cmd
/usr/local/bin/vmtoolsd
/usr/bin/vm-support
/etc/profile.d/open-vm-tools.sh
/etc/vmware-tools/poweron-vm-default
/etc/vmware-tools/scripts/vmware/network
/etc/vmware-tools/poweroff-vm-default
/etc/vmware-tools/resume-vm-default
/etc/vmware-tools/statechange.subr
/etc/vmware-tools/tools.conf.example
/etc/vmware-tools/suspend-vm-default
